# deployment/votercheck.conf
#start on started nginx
#stop on stopped nginx

start on runlevel [2345]
stop on runlevel [016]

exec uwsgi --socket 127.0.0.1:6546 --chdir /home/ubuntu/votercheck/votercheck/web -w wsgi
